# docker-library/django/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/local/uwsgi/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
