# docker-library/django/supervisor.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/site/wwwroot/django/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
