# docker-library/django/supervisor.conf
[program:nginx-app]
command = /usr/sbin/nginx

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/nginx/uwsgi/uwsgi.ini
